# rootfs/etc/supervisord.d/dbus-daemon.ini
[program:dbus-daemon]
command = /bin/bash -c 'mkdir /run/dbus && exec dbus-daemon --system --nopidfile --nofork'
priority = 100
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
