# rootfs/etc/supervisord.d/libvirtd.ini
[program:libvirtd]
command = /sbin/libvirtd
priority = 300
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
